# root/etc/services.d/unifi/run
#!/usr/bin/with-contenv bash

cd /config || exit

exec \
	s6-setuidgid abc java -Xmx1024M -jar /usr/lib/unifi/lib/ace.jar start
