# root/etc/services.d/unifi/run
#!/usr/bin/with-contenv bash

cd /config || exit

exec \
    s6-setuidgid abc java -Dlog4j2.formatMsgNoLookups=true -jar /usr/lib/unifi/lib/ace.jar start;
