# root/etc/services.d/unifi/run
#!/usr/bin/with-contenv bash

cd /config || exit

if [ -z ${MEM_LIMIT+x} ]; then
  MEM_LIMIT="1024M"
fi

exec \
	s6-setuidgid abc java -Xmx"${MEM_LIMIT}" -jar /usr/lib/unifi/lib/ace.jar start
