# services/unifi/run
#!/bin/bash


exec /sbin/setuser abc java -Xmx1024M -jar /usr/lib/unifi/lib/ace.jar start