# services/unifi/run
#!/bin/bash
exec /sbin/setuser abc "/command/to/execute" 